# roles/lamp/templates/hci0_up.conf
[program:hci0_up]
command=/bin/hciconfig hci0 up
priority=200
autostart=true
autorestart=false
startretries=1

