# roles/lamp/templates/lamp_service.conf
[program:lamp_service]
command={{ repo_directory }}/lamp/lamp_service.py
directory={{ repo_directory }}/lamp/
environment=HOME="/home/pi"
priority=300
user=pi
autostart=true
autorestart=true

