# roles/cloud/templates/mqtt-daemon.conf
[program:mqtt-daemon]
command={{ django_dir }}/manage.py mqtt-daemon
directory={{ django_dir }}
priority=400
user=django_user
autostart=true
autorestart=true

